# etc/melody.conf
[program:melody]
command=/opt/melody/melody
directory=/opt/melody
stdout_logfile=/opt/melody/melody_supervisor.out
stderr_logfile=/opt/melody/melody_supervisor.err
autostart=true
autorestart=true
